# sync-kde-and-gtk-places/usr/lib/systemd/user/sync-kde-to-gtk-places.path
[Unit]
Description=Sync from kde places to gtk bookmarks
ConditionUser=!root

[Path]
PathModified=%h/.local/share/user-places.xbel
Unit=sync-kde-to-gtk-places.service

[Install]
WantedBy=default.target
