# sync-kde-and-gtk-places/usr/lib/systemd/user/sync-gtk-to-kde-places.path
[Unit]
Description=Sync from gtk places to KDE bookmarks
ConditionUser=!root

[Path]
SystemCallArchitectures=native
PathModified=%h/.config/gtk-3.0/bookmarks
Unit=sync-gtk-to-kde-places.service

[Install]
WantedBy=default.target
